# monitor/etc/udev/rules.d/99-monitor-hotplug.rules
KERNEL=="card0", SUBSYSTEM=="drm", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/home/carlos/.Xauthority", RUN+="/usr/local/bin/hotplug-monitor.sh"
